# contrib/upstart/cilantro.conf
description "Cilantropypy service"

start on runlevel [2345]
stop on runlevel [!2345]

script
	exec >>/var/log/cilantropy.log 2>&1

	# load defaults if default file exists
	test -f "/etc/default/cilantropy" && . /etc/default/cilantropy

	# makes sure HOST and PORT are set
	test "$HOST" != "" || HOST=127.0.0.1
	test "$PORT" != "" || PORT=5000

	echo starting Cilantropy on $HOST:$PORT
	exec cilantropy --host=$HOST --port=$PORT
end script